# ssl/crontab
0 2 * * 6 /usr/bin/certbot renew --no-self-upgrade
